# openconnect/.config/systemd/user/tuvpn.service
[Unit]
Description=OpenConnect VPN (user service)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/openconnect/vpn.env

ExecStart=/bin/sh -lc 'secret-tool lookup service openconnect user "$VPN_USERNAME" \
  | /usr/bin/openconnect --config "$HOME/.config/openconnect/vpn.conf" \
      --passwd-on-stdin --non-inter \
      --user "$VPN_USERNAME" \
      ${AUTHGROUP:+--authgroup "$AUTHGROUP"} \
      "$VPN_SERVER"'
Restart=on-failure

[Install]
WantedBy=default.target

